# lib/upstart/docker-upstreams.conf
# docker-virtualhosts - Nginx config generator for Docker containers using the NGINX_UPSTREAM env

description "Nginx config generator for Docker containers using the VIRTUAL_HOST env"
author "Romke van der Meulen <redge.online@gmail.com>"

# When to start the service
start on filesystem and started docker

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Send output to logfile
console log

# Run before process
pre-start script
    [ -d /etc/nginx/certs ] || mkdir -p /etc/nginx/certs
end script

# Start the process
exec /etc/nginx/docker-gen-upstream-service
